# qubes-ssh-agent-proxy@.service
[Service]
ExecStart=/usr/bin/qrexec-client-vm "" qubes.SshAgent
StandartOutput=socket
StandardInput=socket
User=user
Group=user